# deploy/com.charles.atb-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.charles.atb-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/charlesow/Desktop/ai-trading-bot/scripts/run_daily.sh</string>
    </array>

    <!-- Weekdays (Mon=1..Fri=5) at 07:30 local — an hour after the 6:30 PT open,
         so the post-earnings setups + marks are available. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/charlesow</string>
        <!-- Set ATB_EXECUTE to 1 here to place paper orders unattended. Leave 0
             (dry run) until you trust a few days of dry output. -->
        <key>ATB_EXECUTE</key>
        <string>0</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/atb-daily.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/atb-daily.err.log</string>
</dict>
</plist>
